# .services/grpc.service
[Unit]
Description=gRPC Server
After=network.target

[Service]
User=a
WorkingDirectory=/home/a/wb_bot_pro/grpc_app
Environment="PYTHONPATH=/home/a/wb_bot_pro/grpc_app"
ExecStart=/home/a/wb_bot_pro/.venv/bin/python /home/a/wb_bot_pro/grpc_app/run.py
EnvironmentFile=/home/a/wb_bot_pro/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target